# root/etc/systemd/system/mnt-nfs.automount
[Unit]
Description=Automount local nfs share

[Automount]
Where=/mnt/nfs
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
